# build_files/files/etc/skel/.config/systemd/user/browsers-1password-setup.service
[Unit]
Description=Setup browsers-1password distrobox
Documentation=man:distrobox-assemble(1)
After=default.target
ConditionPathExists=/etc/distrobox/browsers-1password.ini

[Service]
Type=oneshot
RemainAfterExit=yes
# Check if distrobox doesn't already exist before running (skip if any container exists with this name)
ExecCondition=/bin/bash -c '! distrobox list --no-color 2>/dev/null | grep -q "browsers-1password"'
# Create the distrobox using the configuration file
ExecStart=/usr/bin/distrobox assemble create --file /etc/distrobox/browsers-1password.ini

[Install]
WantedBy=default.target
