# build_files/files/etc/skel/.config/systemd/user/browsers-1password-setup.service
[Unit]
Description=Setup browsers-1password distrobox
Documentation=man:distrobox-assemble(1)
After=default.target
ConditionPathExists=/etc/distrobox/browsers-1password.ini
# CRITICAL: Only run if the container does NOT exist at all
ConditionPathExists=!/var/lib/containers/storage/overlay-containers/

[Service]
Type=oneshot
RemainAfterExit=yes
# Double-check the container doesn't exist before creating
# This prevents accidental deletion of existing containers
ExecStartPre=/bin/bash -c 'if podman container exists browsers-1password; then echo "Container browsers-1password already exists, skipping setup"; exit 1; fi'
# Create the distrobox using the configuration file (WITHOUT replace flag in INI)
ExecStart=/usr/bin/distrobox assemble create --file /etc/distrobox/browsers-1password.ini

[Install]
WantedBy=default.target
